# recipes-core/systemd/files/dashsight.rules
KERNEL=="ttyS2", RUN+="/bin/stty -F /dev/ttyS2 230400 raw -echo"

